# scripts/warm_cache.timer
# /etc/systemd/system/yapit-cache-warm.timer
# Copy to VPS, then: systemctl daemon-reload && systemctl enable --now yapit-cache-warm.timer

[Unit]
Description=Daily cache warming for voice previews

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
